# ignite/etc/sv/acpid/run
#!/bin/sh
[[ -f /etc/conf.d/acpid ]] && . /etc/conf.d/acpid
exec acpid -f $ACPID_ARGS
