# ignite/etc/sv/dhcpcd/run
#!/bin/bash
exec dhcpcd -B -d -D -H -L -t 0
