# ignite/etc/sv/cpupower/run
#!/bin/sh
set -e
/usr/lib/systemd/scripts/cpupower
exec chpst -b cpupower pause
