# ignite/etc/sv/nfs-server/finish
#!/bin/sh
exec 2>&1
# shut down the nfsd threads.
/usr/sbin/rpc.nfsd -- 0
/usr/sbin/exportfs -a -u
