# ignite/etc/sv/shorewall6/finish
#!/bin/sh
exec /usr/sbin/shorewall6 stop 2>&1
