# ignite/etc/sv/ntpd/run
#!/bin/sh
exec ntpd -g -u ntp:ntp -n 2>&1
