# ignite/etc/sv/mysqld/run
#!/bin/sh

exec /usr/bin/mysqld --user=mysql 2>&1
