# ignite/etc/sv/nginx/run
#!/bin/sh

/usr/sbin/nginx -t -q -c /etc/nginx/nginx.conf
test $? -ne 0 ||
  exec /usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon off;'

echo "nginx config test failed."
sv down .
