# ignite/etc/sv/sshd/run
#!/bin/sh

/usr/bin/ssh-keygen -A # Will generate host keys if they don't already exist
exec /usr/bin/sshd -D
