# ignite/etc/sv/openvpn/run
#!/bin/sh

exec /usr/sbin/openvpn --cd /etc/openvpn --config openvpn.conf
