# ignite/etc/sv/apache/run
#!/bin/bash -e
mkdir -p /run/httpd
[ -r /usr/bin/envvars ] && . /usr/bin/envvars
[ -r /etc/conf.d/apache ] && . /etc/conf.d/apache
exec httpd -DNO_DETACH 2>&1
