# ignite/etc/sv/mpdscribble/run
#!/bin/sh
exec /usr/bin/mpdscribble --no-daemon
