# ignite/etc/sv/eth0/run
#!/bin/sh
ifconfig eth0 up
exec dhcpcd -qLB -t 0 eth0
