# ignite/etc/sv/eth0/run
#!/bin/sh
exec dhcpcd -qLB -t 0 -C wpa_supplicant -C timezone -C hostname ${PWD##*/}
