# ignite/sv/test/run
#!/bin/sh
echo stdout
echo stderr >/dev/stderr

sleep 10
