# ignite/etc/sv/laptop-mode/finish
#!/bin/sh
rm -f /var/run/laptop-mode-tools/enabled
exec /usr/sbin/laptop_mode stop >/dev/null 2>&1
