# ignite/etc/sv/laptop-mode/run
#!/bin/sh
set -e
mkdir -p /var/run/laptop-mode-tools
touch /var/run/laptop-mode-tools/enabled
/usr/sbin/laptop_mode auto >/dev/null 2>&1
exec chpst -b laptop-mode pause
