# ignite/etc/sv/shorewall/finish
#!/bin/sh
exec /usr/sbin/shorewall stop 2>&1
