# ignite/etc/sv/ddclient/run
#!/bin/bash

[ -f /etc/conf.d/ddclient ] && . /etc/conf.d/ddclient

exec /usr/sbin/ddclient -foreground -syslog ${EXTRA_ARGS}
