# ignite/etc/sv/dhcp4/run
#!/bin/sh

. /etc/conf.d/dhcp

/usr/sbin/dhcpd -4 -p -q -pf "/run/dhcpd.pid" $DHCP4_ARGS;
