# ignite/etc/sv/mpd/run
#!/bin/sh
exec /usr/bin/mpd --no-daemon /etc/mpd.conf
