# ignite/etc/sv/agetty-ttyS0/run
#!/bin/sh
service=$(basename $(pwd))
tty=${service/*-}

exec setsid /sbin/agetty -8 -s 9600 --noclear $tty linux
