# ignite/etc/sv/iptables/finish
#!/bin/bash
/usr/lib/systemd/scripts/iptables-flush
