# ignite/etc/sv/iptables/run
#!/bin/sh
set -e
/usr/bin/iptables-restore /etc/iptables/iptables.rules
exec chpst -b iptables pause
