# ignite/etc/sv/tftpd/run
#!/bin/sh
. /etc/conf.d/tftpd
exec /usr/sbin/in.tftpd -L $TFTPD_ARGS
