# ignite/etc/sv/net-auto-wired/finish
#!/bin/sh
. /etc/conf.d/netcfg
exec /usr/bin/netcfg iface-down $WIRED_INTERFACE
