# ignite/etc/sv/avahi-daemon/run
#!/bin/sh
sv check dbus >/dev/null &&
exec avahi-daemon -s
