# ignite/etc/sv/dhcp6/run
#!/bin/sh

. /etc/conf.d/dhcp

/usr/sbin/dhcpd -6 -f -q -pf "/run/dhcpd.pid" $DHCP6_ARGS;
