# ignite/etc/sv/dhcp6/run
#!/bin/sh
. /etc/conf.d/dhcp
exec /usr/sbin/dhcpd -6 -f --no-pid $DHCP6_ARGS
