# ignite/etc/sv/haveged/run
#!/bin/sh
exec /usr/sbin/haveged -F -w 1024 -v 1
