# ignite/etc/sv/rpc.statd/run
#!/bin/sh -e
. /etc/conf.d/nfs-common.conf
sv check rpcbind >/dev/null
exec /usr/bin/rpc.statd -F $STATD_OPTS
