# ignite/etc/sv/gpm/finish
#!/bin/sh
/usr/sbin/gpm -D -k 2>&1 | logger -s -p daemon.notice -t gpm
