# ignite/etc/sv/gpm/run
#!/bin/sh
/usr/sbin/gpm -D -m /dev/input/mice -t imps2 2>&1 | logger -s -p daemon.notice -t gpm
