# ignite/etc/sv/cupsd/run
#!/bin/sh
sv check colord >/dev/null && sv check colord-sane >/dev/null &&
exec cupsd -f
