# ignite/etc/sv/rpcbind/finish
#!/bin/sh
umount /var/lib/nfs/rpc_pipefs || true
