# ignite/etc/sv/rpcbind/run
#!/bin/sh -e
mountpoint -q /var/lib/nfs/rpc_pipefs ||
  mount -t rpc_pipefs rpc_pipefs /var/lib/nfs/rpc_pipefs
exec /usr/bin/rpcbind -f
