# ignite/etc/sv/wlan0-wpa/run
#!/bin/sh
rm -f /run/wpa_supplicant/$(cat interface)
exec wpa_supplicant -i $(cat interface) -D nl80211,wext -c/etc/wpa_supplicant/wpa_supplicant.conf 2>&1
