# ignite/etc/sv/net-auto-wireless/finish
#!/bin/sh
. /etc/conf.d/netcfg
exec /usr/bin/netcfg-wpa_actiond stop $WIRELESS_INTERFACE
