# ignite/etc/sv/net-auto-wireless/run
#!/bin/sh
. /etc/conf.d/netcfg
/usr/bin/netcfg-wpa_actiond $WIRELESS_INTERFACE
exec chpst -b net-auto-wireless pause
