# ignite/etc/sv/wicd/run
#!/bin/sh -e

echo "Starting wicd Daemon"
echo "Killing all instances of wicd-daemon.py"
pkill -f wicd-daemon.py

/usr/sbin/wicd --no-daemon
