# ignite/etc/sv/wicd/run
#!/bin/sh
sv check dbus >/dev/null &&
exec wicd --no-daemon
