# ignite/etc/sv/hwclock/finish
#!/bin/bash
. /etc/rc.conf
hwclock --adjust ${HARDWARECLOCK:+--$(echo $HARDWARECLOCK |tr A-Z a-z)}
