# ignite/etc/sv/rpc.idmapd/run
#!/bin/sh
. /etc/conf.d/nfs-common.conf
sv check nfsd >/dev/null
exec /usr/bin/rpc.idmapd -f $IDMAPD_OPTS 2>&1
