# ignite/etc/sv/dnsmasq/run
#!/bin/sh
exec /usr/bin/dnsmasq -k --user=dnsmasq
