# ignite/etc/sv/statd/run
#!/bin/sh

exec 2>&1

# Make sure the portmap service is running.
sv check rpcbind

exec rpc.statd -F -d
