# ignite/etc/sv/dkms/finish
#!/bin/sh
exec /usr/lib/systemd/scripts/dkms stop
