# ignite/etc/sv/dkms/run
#!/bin/sh
set -e
/usr/lib/systemd/scripts/dkms start >/dev/null
exec chpst -b dkms pause
