# ignite/etc/sv/preload/run
#!/bin/sh
. /etc/conf.d/preload
exec /usr/bin/preload --foreground $PRELOAD_OPTS
