# ignite/etc/sv/nfsd/finish
#!/bin/sh
/usr/bin/exportfs -a -u
/usr/bin/rpc.nfsd 0
