# ignite/etc/sv/wlan0/run
#!/bin/sh
exec dhcpcd -qLB -t 0 ${PWD##*/}
